# templates/pgbouncer-follow-primary.timer.erb
[Unit]
Description=Run pgbouncer-follow-primary periodically

[Timer]
OnBootSec=10s
OnUnitActiveSec=<%= interval %>s
AccuracySec=1s
Unit=pgbouncer-follow-primary.service

[Install]
WantedBy=timers.target
